# mkosi.extra/etc/systemd/system/keepalived.service.d/override.conf
[Unit]
Requires=network-online.target
After=network-online.target

[Service]
PIDFile=/run/keepalived.pid
KillMode=control-group
